# backend/brackett-api.service
[Unit]
Description=brackett.dev FastAPI backend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/brackett.dev/backend
Environment="PATH=/var/www/brackett.dev/backend/.venv/bin"
EnvironmentFile=/var/www/brackett.dev/backend/.env.production
ExecStart=/var/www/brackett.dev/backend/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
